# base/keycfg/files/keycfg.init
#!/bin/sh /etc/rc.common
START=20

config_cb() {
	local cfg="$CONFIG_SECTION"
	config_get configname "$cfg" TYPE

	case "$configname" in
		authorized_key)
			config_get data "$cfg" data
			echo "$data" >> /etc/dropbear/authorized_keys			
		;;
	esac
}

start() {
	umask 077
	mkdir -p /etc/dropbear
	rm -f /etc/dropbear/authorized_keys
	config_load keycfg
}

